# staticmirror/etc/cron.daily/staticmirror-scheduler.sh
#!/bin/bash
#run main script for mirror-static
/etc/mirror-static/main.sh


